# root/etc/s6-overlay/s6-rc.d/svc-ai-chat-app/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/ai-chat-app

exec s6-setuidgid abc \
    yarn start
