# root/etc/s6-overlay/s6-rc.d/init-config-ai-chat-app/run
#!/usr/bin/with-contenv bash

# initialise database
cd /app/ai-chat-app
yarn prisma migrate deploy

# permissions
find /app/ai-chat-app -path "*/node_modules" -prune -o -exec chown abc:abc {} +
chown -R abc:abc \
    /config \
    /root
