# root/etc/s6-overlay/s6-rc.d/init-config-ai-chat-app/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# initialise database
cd /app/ai-chat-app
npx prisma migrate deploy &>/dev/null

# permissions
chown -R abc:abc \
    /config \
    /app
